# apache2-ipv6local/etc/systemd/system/apache2.service.d/update-ipv6local.conf
# Runs update upv6local once before apache starts.
# The script runs once to update /etc/apache2/conf-available/00-ipv6.conf
[Unit]
Requires=update-ipv6local.service
After=update-ipv6local.service

