# apache2-ipv6local/etc/systemd/system/update-ipv6local.service
[Unit]
Description=Update IPv6Local Prefix

[Service]
Type=oneshot
ExecStart=/usr/bin/update-ipv6local.sh

